# deploy/systemd/eveai-codex-proxy.service
[Unit]
Description=EVEAI Codex OpenAI Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Current production layout keeps auth under /root/.codex/auth/.
# If auth is moved to a dedicated service account, update User and auth-path together.
User=root
WorkingDirectory=/opt/codex_proxy
Environment=RUST_LOG=info
ExecStart=/opt/codex_proxy/target/release/codex-openai-proxy --port 8080 --auth-path /root/.codex/auth/
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
