# deploy/systemd/eveai.service
[Unit]
Description=EVE AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust these paths for your server before installing.
User=eveai
Group=eveai
WorkingDirectory=/srv/eveai
EnvironmentFile=/srv/eveai/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /srv/eveai/dist/app.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# The checkout, package metadata, and built code remain read-only. Runtime
# state and the adjacent single-process lock live under this explicit path.
ReadWritePaths=/srv/eveai/data

[Install]
WantedBy=multi-user.target
