# deploy/systemd/eveai-backend.service
[Unit]
Description=EVEAI Backend
After=network-online.target eveai-codex-proxy.service
Wants=network-online.target eveai-codex-proxy.service

[Service]
Type=simple
WorkingDirectory=/home/antipedik/eveai
EnvironmentFile=/home/antipedik/eveai/.env
ExecStart=/home/antipedik/.nvm/versions/node/v22.22.0/bin/node /home/antipedik/eveai/dist/app.js
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
